# sysutils/dockerbox/pkg-plist
etc/dockerbox/dockerbox.conf
etc/rc.d/dockerbox
sbin/dockerbox
