# linux/dockerbox/pkg-plist
etc/dockerbox/dockerbox.conf
etc/rc.d/dockerbox
sbin/dockerbox
share/dockerbox/device.map
share/dockerbox/disk.img
